# systemd/voxnode-recorder.service
[Unit]
Description=voxnode audio recorder (continuous Opus segments)
Documentation=https://github.com/dialytica/voxnode
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
User=contai
Group=contai
# Supplementary groups — нужны для доступа к /dev/snd (группа audio)
SupplementaryGroups=audio

# Каталог установки и venv (ohmyzsh-style: всё в /opt/voxnode)
Environment=VOXNODE_HOME=/opt/voxnode
Environment=VOXNODE_CONFIG=/etc/voxnode/config.yaml
Environment=PYTHONUNBUFFERED=1
# HOME и LANG обязательны для ALSA — без них snd_config_get_card не находит карту по имени
Environment=HOME=/home/contai
Environment=LANG=en_US.UTF-8

ExecStart=/opt/voxnode/.venv/bin/python -m voxnode.recorder

WorkingDirectory=/opt/voxnode

# Restart: ffmpeg должен крутиться 24/7. При падении — мгновенный рестарт.
Restart=always
RestartSec=3

# Логи в journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voxnode-recorder

# Hardening: минимальные привилии.
# ВАЖНО: не используем ProtectSystem/DevicePolicy — ALSA требует доступ к
# /proc/asound, /dev/snd/*, и shared memory. Они ломают запись.
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target
